# systemd/kiosk.service
[Unit]
Description=Chromium Kiosk
Requires=signage.service
After=signage.service

[Service]
User=pi
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=DISPLAY=:0
ExecStart=/usr/bin/chromium --no-sandbox --kiosk --incognito http://localhost:8000
Restart=always

[Install]
WantedBy=multi-user.target
